# thermpro-monitor.service
[Unit]
Description=ThermoPro BLE Hygrometer Monitor
After=bluetooth.service
Wants=bluetooth.target

[Service]
Type=simple
User=pi
Group=pi
SupplementaryGroups=bluetooth
WorkingDirectory=/opt/thermpro-monitor
EnvironmentFile=/opt/thermpro-monitor/thermpro-monitor.env
ExecStart=/opt/thermpro-monitor/.venv/bin/python /opt/thermpro-monitor/thermpro_monitor.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
